# tools/docker/edge/ui/root/etc/s6-overlay/s6-rc.d/init-openems/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cat << EOF > /var/www/html/openems/appsettings.json
{
	"edgeWebsocket": http://${EDGE_HOST:-127.0.0.1}:${EDGE_WEBSOCKET_PORT:-8085};
}
EOF